# packages/zenlink-mcp/scripts/launchd-zenlink-mcp-daemon.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Example launchd plist when zenlink-mcp --daemon is started OUTSIDE OpenClaw's own MCP
  supervisor (e.g. you wrap `node …/dist/cli.js --daemon` yourself).

  Paths: use ABSOLUTE paths everywhere. Do not point ProgramArguments at a versioned extract
  under /tmp/zenlink-mcp-openclaw-macos-v0.12.x/ — that path changes every release. Prefer a stable
  install, e.g. ~/.openclaw/zenlink-mcp/current/zenlink-mcp/dist/cli.js (see upgrade-offline-install.sh
  in the offline tarball). On macOS, /tmp is a symlink to /private/tmp; both resolve to the same
  files, but launchd log paths must still be real paths your user can write.

  launchd does NOT source zenlink-deploy.env. Either:
  - mirror the same key/value pairs from ~/.openclaw/openclaw.json → mcp.servers.<name>.env, or
  - list them here under EnvironmentVariables (keep token secret; chmod 600 this file).

  ZENLINK_MCP_USE_DAEMON on THIS job: keep 0 — this process *is* the daemon (--daemon). The stdio
  MCP workers (spawned by OpenClaw) use USE_DAEMON=1 to forward to the addr file below.

  If OpenClaw spawns stdio MCP / daemon from mcp.servers env, you usually do NOT need this plist
  for hook push — fix registration + restart Gateway first.
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>net.zenheart.zenlink-mcp-daemon</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/REPLACE/WITH/ABSOLUTE/PATH/TO/stable/zenlink-mcp/dist/cli.js</string>
    <string>--daemon</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ZENLINK_AGENT_ID</key>
    <string></string>
    <key>ZENLINK_TOKEN</key>
    <string></string>
    <key>ZENLINK_MCP_OPENCLAW_HOOK_BASE</key>
    <string>http://127.0.0.1:18789/hooks</string>
    <key>ZENLINK_MCP_OPENCLAW_HOOK_TOKEN</key>
    <string>same-as-openclaw.json-hooks-token</string>
    <key>ZENLINK_MCP_OPENCLAW_WAKE_MODE</key>
    <string>now</string>
    <!-- Optional request-level wake routing; installer normally sets hooks.defaultSessionKey instead. -->
    <!-- <key>ZENLINK_MCP_OPENCLAW_SESSION_KEY</key>
    <string>hook:zenheart-main</string> -->
    <key>ZENLINK_MCP_DAEMON_ADDR_FILE</key>
    <string>/REPLACE/WITH/ABSOLUTE/PATH/.openclaw/tmp/zenlink-mcp-daemon.addr</string>
    <key>ZENLINK_MCP_USE_DAEMON</key>
    <string>0</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/.openclaw/tmp/zenlink-mcp-daemon.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/REPLACE/WITH/ABSOLUTE/PATH/.openclaw/tmp/zenlink-mcp-daemon.launchd.log</string>
</dict>
</plist>
